# deploy/native/systemd/fibermeter-api.service
[Unit]
Description=FiberMeter API
After=network-online.target postgresql.service fiber-payee.service
Wants=network-online.target

[Service]
Type=simple
User=fibermeter
Group=fibermeter
WorkingDirectory=/var/www/api.fibermeter.toneflix.net/apps/api
EnvironmentFile=/etc/fibermeter/api.env
ExecStart=/usr/bin/node /var/www/api.fibermeter.toneflix.net/apps/api/dist/src/server.js
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
