# deploy/native/systemd/fibermeter-autopay.service
[Unit]
Description=FiberMeter bounded testnet auto-payer
After=network-online.target fibermeter-api.service fiber-payer.service
Requires=fibermeter-api.service fiber-payer.service

[Service]
Type=simple
User=fibermeter
Group=fibermeter
WorkingDirectory=__APP_ROOT__
EnvironmentFile=/etc/fibermeter/autopay.env
ExecStart=/usr/bin/node __APP_ROOT__/scripts/autopay.mjs
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
